# systemd/daily-pnl.timer
[Unit]
Description=Run Hyperliquid Daily PnL report at 07:05 Europe/Berlin

[Timer]
OnCalendar=*-*-* 07:05:00
Timezone=Europe/Berlin
Unit=daily-pnl.service
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target

