# systemd/midday-pnl.service
[Unit]
Description=Hyperliquid MM - Midday PnL Report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jerry/hyperliquid-mm-bot-complete
User=jerry
EnvironmentFile=-/home/jerry/hyperliquid-mm-bot-complete/.env
Environment="NODE_ENV=production"
Environment="PNL_PROFILE=MIDDAY"
Environment="PNL_SKIP_HISTORY=true"
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/jerry/hyperliquid-mm-bot-complete
PrivateTmp=true
ExecStart=/usr/bin/env bash /home/jerry/hyperliquid-mm-bot-complete/scripts/daily-pnl.sh
Restart=no

[Install]
WantedBy=multi-user.target
