# systemd/mm-bot-health.timer
[Unit]
Description=Run mm-bot health watchdog every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
RandomizedDelaySec=5
Unit=mm-bot-health.service
Persistent=true

[Install]
WantedBy=timers.target



